# .config/uwsgi/uwsgi.service
[Unit]
Description=uWSGI service
After=syslog.target

[Service]
ExecStart=/home/jeon/project/tactic-on-table-django-backend/TacticOnTableVenv/bin/uwsgi -i /home/jeon/project/tactic-on-table-django-backend/.config/uwsgi/tactic_on_table.ini

Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target